# infra/sing-box/config/init/runit/sing-box/run
#!/bin/sh
exec 2>&1
exec /usr/local/bin/sing-box run -D /etc/sing-box
